# train_switch.service
[Unit]
Description=Train Switch
# launch service after network is established
After=network.target

[Service]
ExecStart=/usr/bin/python3 app.py
# ensure to replace `User` with your user
WorkingDirectory=/home/celery/Pi-Train-Switching
StandardOutput=inherit
StandardError=inherit
Restart=always
# ensure to replace `User` with your user
User=celery

[Install]
WantedBy=multi-user.target
